# frontend/deploy/raspberry-pi/pcb-react-kiosk-web.service
[Unit]
Description=PCB React Kiosk Web App (Vite Preview)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/AiCapstoneV2/frontend
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 4173
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
