# frontend/deploy/raspberry-pi/pcb-react-kiosk-web.service
[Unit]
Description=PCB React Kiosk Web App (Vite Preview)
After=network-online.target pcb-edge.service
Wants=network-online.target pcb-edge.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/AiCapstoneV2/frontend
Environment=NODE_ENV=production
EnvironmentFile=/etc/default/aicapstone
Environment=VITE_API_PROXY_TARGET=http://${API_HOST}:8080
Environment=VITE_EDGE_CAPTURE_URL=http://${EDGE_HOST}:8000
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 4173
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
